# diaspora.target
[Unit]
Description=Diaspora social network
StopWhenUnneeded=true

[Install]
WantedBy=multi-user.target
