# daemontools/worker/run
#!/bin/bash

cd  "$(dirname $(readlink -e $0))/.."
source ./common.sh

export QUEUE='*'

exec setuidgid "$DIASPORA_USER" bundle exec rake resque:work